# systemd/openclaw-memory-weekly.timer
[Unit]
Description=Weekly full memory curation every Sunday at 03:00 AM

[Timer]
OnCalendar=Sun *-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
