# systemd/openclaw-memory-daily.timer
[Unit]
Description=Daily memory quick tidy at 06:00 AM

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
