# systemd/openclaw-security-weekly.timer
[Unit]
Description=Weekly full security audit every Saturday at 04:00 AM

[Timer]
OnCalendar=Sat *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
